# scripts/creckstars-backup.timer
# Install on the VM: /etc/systemd/system/creckstars-backup.timer
# Enable: sudo systemctl daemon-reload && sudo systemctl enable --now creckstars-backup.timer
# Check:  systemctl list-timers creckstars-backup
[Unit]
Description=Run CrecKStars backup nightly

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true        # if the VM was off at 03:00, run on next boot
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
